# kime-xim.service
[Unit]
Description=Kime XIM Daemon

[Service]
Type=simple
ExecStart=/usr/bin/kime-xim
Restart=on-failure

[Install]
WantedBy=multi-user.target

